# systemd/claws-mail.service
[Unit]
Description=Claws-Mail

[Service]
Type=simple
ExecStart=/bin/claws-mail --online
Restart=always

[Install]
WantedBy=multi-user.target
