# systemd/mail.service
[Unit]
Description=Claws-Mail

[Service]
Type=simple
ExecStart=/home/luca-linux/bin/geary_start
Restart=always

[Install]
WantedBy=multi-user.target
