# systemd/touchegg-user.service
[Unit]
Description=Touchégg
After=touchegg.service

[Service]
ExecStart=/usr/bin/touchegg --client
Restart=always

[Install]
WantedBy=default.target
