# _ansible/roles/distro_setup/files/systemd/user/tmpfiles-cleanup.timer
[Unit]
Description=Tmpfiles Automatic Cleanup Trigger
Documentation=man:find(1)

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1w

[Install]
WantedBy=timers.target
